# prometheus/ansible_config/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
[Service]
User={{ node_exporter_user }}
ExecStart={{ install_base }}/bin/node_exporter
[Install]
WantedBy=default.target
